# platform-diagnostics/misc/osprey-diagnostics.service
[Unit]
Description=Osprey diagnostics platform

[Service]
Type=oneshot
ExecStart=/opt/osprey/diagnostics/platform-diagnostics
User=pi